# installer/debian/DEBIAN/postinst
#!/bin/bash
# Post-installation script for canvuslocallm
set -e

INSTALL_DIR="/opt/canvuslocallm"
SERVICE_USER="canvusllm"
SERVICE_GROUP="canvusllm"
SERVICE_NAME="canvuslocallm"
LIB_DIR="$INSTALL_DIR/lib"
MODELS_DIR="$INSTALL_DIR/models"

case "$1" in
    configure)
        # Create installation directory if not exists
        if [ ! -d "$INSTALL_DIR" ]; then
            mkdir -p "$INSTALL_DIR"
            echo "Created installation directory: $INSTALL_DIR"
        fi

        # Create lib directory for shared libraries
        if [ ! -d "$LIB_DIR" ]; then
            mkdir -p "$LIB_DIR"
            echo "Created lib directory: $LIB_DIR"
        fi

        # Create models directory for AI models
        if [ ! -d "$MODELS_DIR" ]; then
            mkdir -p "$MODELS_DIR"
            echo "Created models directory: $MODELS_DIR"
        fi

        # Create service group if not exists
        if ! getent group "$SERVICE_GROUP" > /dev/null 2>&1; then
            groupadd --system "$SERVICE_GROUP"
            echo "Created system group: $SERVICE_GROUP"
        fi

        # Create service user if not exists
        if ! getent passwd "$SERVICE_USER" > /dev/null 2>&1; then
            useradd --system \
                --gid "$SERVICE_GROUP" \
                --home-dir "$INSTALL_DIR" \
                --shell /usr/sbin/nologin \
                --comment "CanvusLocalLLM Service Account" \
                "$SERVICE_USER"
            echo "Created system user: $SERVICE_USER"
        fi

        # Add service user to video group for GPU access
        if getent group video > /dev/null 2>&1; then
            usermod -a -G video "$SERVICE_USER" 2>/dev/null || true
            echo "Added $SERVICE_USER to video group for GPU access"
        fi

        # Copy .env.example to .env if .env doesn't exist
        if [ ! -f "$INSTALL_DIR/.env" ]; then
            if [ -f "$INSTALL_DIR/.env.example" ]; then
                cp "$INSTALL_DIR/.env.example" "$INSTALL_DIR/.env"
                echo "Created configuration file: $INSTALL_DIR/.env"
                echo "Please edit $INSTALL_DIR/.env with your settings"
            fi
        fi

        # Set ownership and permissions
        chown -R "$SERVICE_USER:$SERVICE_GROUP" "$INSTALL_DIR"
        chmod 750 "$INSTALL_DIR"

        # Protect configuration file if it exists
        if [ -f "$INSTALL_DIR/.env" ]; then
            chmod 640 "$INSTALL_DIR/.env"
        fi

        # Make binary executable if it exists
        if [ -f "$INSTALL_DIR/canvuslocallm" ]; then
            chmod 755 "$INSTALL_DIR/canvuslocallm"
        fi

        # Set library permissions
        if [ -d "$LIB_DIR" ]; then
            chmod 755 "$LIB_DIR"
            # Make shared libraries executable/readable
            find "$LIB_DIR" -name "*.so*" -exec chmod 755 {} \; 2>/dev/null || true
        fi

        # Set model file permissions (read-only for service)
        if [ -d "$MODELS_DIR" ]; then
            chmod 755 "$MODELS_DIR"
            find "$MODELS_DIR" -type f -exec chmod 644 {} \; 2>/dev/null || true
        fi

        # Configure LD_LIBRARY_PATH for the shared library
        # Create ld.so.conf.d entry for canvuslocallm libs
        LD_CONF_FILE="/etc/ld.so.conf.d/canvuslocallm.conf"
        if [ ! -f "$LD_CONF_FILE" ]; then
            echo "$LIB_DIR" > "$LD_CONF_FILE"
            echo "Created library config: $LD_CONF_FILE"
        fi

        # Update dynamic linker cache
        ldconfig

        # Install systemd service
        if [ -f "$INSTALL_DIR/canvuslocallm.service" ]; then
            cp "$INSTALL_DIR/canvuslocallm.service" /etc/systemd/system/
            echo "Installed systemd service"
        fi

        # Reload systemd daemon
        systemctl daemon-reload

        # Enable the service
        systemctl enable "$SERVICE_NAME.service"
        echo "Enabled $SERVICE_NAME service"

        # Check for NVIDIA GPU
        GPU_DETECTED=false
        if command -v nvidia-smi > /dev/null 2>&1; then
            if nvidia-smi > /dev/null 2>&1; then
                GPU_DETECTED=true
                GPU_INFO=$(nvidia-smi --query-gpu=name,memory.total --format=csv,noheader 2>/dev/null | head -1)
                echo "Detected NVIDIA GPU: $GPU_INFO"
            fi
        fi

        # Print post-installation instructions
        echo ""
        echo "========================================"
        echo "CanvusLocalLLM Installation Complete"
        echo "========================================"
        echo ""
        echo "Installed components:"
        echo "  - Main application: $INSTALL_DIR/canvuslocallm"
        if [ -f "$LIB_DIR/libstable-diffusion.so" ]; then
            echo "  - SD Library: $LIB_DIR/libstable-diffusion.so"
        fi
        if [ -f "$MODELS_DIR/sd-v1-5.safetensors" ]; then
            echo "  - SD Model: $MODELS_DIR/sd-v1-5.safetensors"
        fi
        echo ""
        if [ "$GPU_DETECTED" = true ]; then
            echo "GPU Status: NVIDIA GPU detected - image generation enabled"
        else
            echo "GPU Status: No NVIDIA GPU detected"
            echo "  Image generation will use CPU (slower) or cloud fallback"
            echo "  To enable GPU acceleration:"
            echo "    sudo apt install nvidia-cuda-toolkit nvidia-driver-535"
        fi
        echo ""
        echo "Next steps:"
        echo "1. Edit configuration: sudo nano $INSTALL_DIR/.env"
        echo "2. Start the service:  sudo systemctl start $SERVICE_NAME"
        echo "3. Check status:       sudo systemctl status $SERVICE_NAME"
        echo "4. View logs:          sudo journalctl -u $SERVICE_NAME -f"
        echo ""
        echo "For more information, see: $INSTALL_DIR/README.txt"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        # Nothing to do
        ;;

    *)
        echo "postinst called with unknown argument: $1" >&2
        exit 1
        ;;
esac

exit 0
